# xdg/autostart/tazpkg-notify.desktop
[Desktop Entry]
Type=Application
Name=Packages Notification System
Name[fr]=Notification du système de paquets
Exec=tazpkg-notify
